# build/files/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/httpd -DFOREGROUND

[program:ssh]
command=/usr/sbin/sshd -D
