# build/files/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/httpd -c "ErrorLog /dev/stdout" -DFOREGROUND
redirect_stderr=true

[program:ssh]
command=/usr/sbin/sshd -D
redirect_stderr=true
